# com.user.wifi.on.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.wifi.on</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/toggle-wifi.sh</string>
        <string>on</string>
    </array>

    <!-- Friday 22:00 (10pm) and Monday 00:00 (midnight starting Monday).
         Weekday values: Sun=0, Mon=1, Tue=2, Wed=3, Thu=4, Fri=5, Sat=6 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>22</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>0</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/wifi-schedule.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/wifi-schedule.log</string>
</dict>
</plist>
